# preferences/etc/cron.d/auto_update

# Run the update script every day at 2 AM
0 2 * * * /home/francky/docker/update_script.sh >> /home/francky/docker/update.log 2>&1
